# apps/server/imwallet.service
[Unit]
Description=IMWallet Server - 区块链钱包后端服务
After=network.target postgresql.service
Wants=postgresql.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=imwallet
Group=imwallet
WorkingDirectory=/opt/imwallet-server
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

# 环境变量文件
EnvironmentFile=/opt/imwallet-server/.env.production

# 安全限制
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/imwallet-server /var/log/imwallet
PrivateTmp=true

# 允许读取密钥文件目录
ReadOnlyPaths=/opt/imwallet-server/keys

# 日志: 同时写入 journal 和文件
StandardOutput=journal+file:/var/log/imwallet/imwallet.log
StandardError=journal+file:/var/log/imwallet/imwallet-error.log
SyslogIdentifier=imwallet

[Install]
WantedBy=multi-user.target